# etc/s6-overlay/s6-rc.d/git-daemon/run
#!/command/execlineb -P
/command/fdmove -c 2 1
/command/s6-setuidgid git
/usr/libexec/git-core/git-daemon --reuseaddr --base-path=/var/lib/git/repositories
