# etc/s6-overlay/s6-rc.d/git-daemon/run
#!/command/execlineb -P
fdmove -c 2 1
/usr/libexec/git-core/git-daemon --user=git --group=git --base-path=/var/lib/git/repositories --reuseaddr
