# etc/s6-overlay/s6-rc.d/openssh-server/run
#!/command/execlineb -P
/command/fdmove -c 2 1
/command/s6-setuidgid git
/usr/sbin/sshd -D -e -p 8022
