# etc/s6-overlay/s6-rc.d/openssh-server/run
#!/command/execlineb -P
s6-setuidgid git
/usr/sbin/sshd -D -e -p 8022
